# lxd-init
#!/sbin/openrc-run
#
# OpenRC service script for LXD
#

command="/usr/sbin/lxd"
command_args="${LXD_OPTIONS}"
command_background="true"
command_user="root:lxd"
pidfile="/var/run/${RC_SVCNAME}.pid"

depend() {
        need net
        use  lxcfs
        after firewall
        provide lxd
}
